# docs/operations/systemd/nzbg.service
[Unit]
Description=nzbg Usenet downloader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nzbg --foreground --config /etc/nzbg.conf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
